# scripts/systemd/encar-monitor.service
[Unit]
Description=Encar Monitor Service
Documentation=https://github.com/lym-afla/Encar
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=user1
Group=user1
WorkingDirectory=/opt/encar-monitor
Environment=PATH=/opt/encar-monitor/venv/bin
EnvironmentFile=/opt/encar-monitor/.env

# Main service command
ExecStart=/opt/encar-monitor/venv/bin/python encar_monitor_api.py

# Service management
Restart=always
RestartSec=10
StartLimitBurst=3
StartLimitInterval=60

# Resource limits
TimeoutStartSec=120
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=encar-monitor

# Playwright environment for systemd services
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/encar-monitor/ms-playwright
Environment=XDG_CACHE_HOME=/opt/encar-monitor/.cache
Environment=PLAYWRIGHT_SKIP_VALIDATE_HOST_REQUIREMENTS=true

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/encar-monitor/logs /opt/encar-monitor/data /opt/encar-monitor/ms-playwright /opt/encar-monitor/.cache

[Install]
WantedBy=multi-user.target
